# CORE/lombix-base/POSTLOAD/usr/etc/profile.d/devel.sh
pathappend /usr/devel/bin		PATH
pathappend /usr/devel/lib/pkgconfig	PKG_CONFIG_PATH
pathappend /usr/devel/share/aclocal	ACLOCAL_PATH
pathappend /usr/devel/share/man		MANPATH
