# CORE/lombix-base/POSTLOAD/usr/etc/profile.d/emulation.sh
pathappend /usr/emulation/bin			PATH
pathappend /usr/emulation/lib/pkgconfig		PKG_CONFIG_PATH
pathappend /usr/emulation/share/aclocal		ACLOCAL_PATH
pathappend /usr/emulation/share/man		MANPATH
