# CORE/lombix-base/POSTLOAD/usr/etc/profile.d/go.sh
pathappend /usr/go/bin PATH
export GOPATH=$HOME/go
export GOROOT=/usr/go
