# CORE/lombix-base/POSTLOAD/usr/etc/profile.d/perl.sh
pathappend /usr/perl/bin			PATH
pathappend /usr/perl/share/man			MANPATH

export PERL=/usr/perl/bin/perl
