# CORE/lombix-base/POSTLOAD/usr/etc/profile.d/perl.sh
pathappend /usr/perl/bin		PATH
pathappend /usr/perl/share/aclocal	ACLOCAL_PATH
pathappend /usr/perl/share/man		MANPATH
