# CORE/lombix-base/POSTLOAD/usr/etc/profile.d/python.sh
pathappend /usr/python/bin			PATH
pathappend /usr/python/share/man		MANPATH
pathappend /usr/python/lib/pkgconfig		PKG_CONFIG_PATH
export PYTHON3=/usr/python/bin/python3
export PYTHON2=/usr/python/bin/python2
