# CORE/lombix-base/POSTLOAD/usr/etc/profile.d/qt5.sh
pathappend /usr/qt5/bin			PATH
pathappend /usr/qt5/lib/pkgconfig	PKG_CONFIG_PATH
pathappend /usr/qt5/share/man			MANPATH

